# files/regulateur_freq.service
[Unit]
Description=Régulateur de fréquence CPU

[Service]
Type=simple
ExecStart=/usr/bin/regulateur_freq.sh

[Install]
WantedBy=multi-user.target
